# MEDSTARGENX/deployment/ml_api.service
[Unit]
Description=MEDSTARGENX ML API Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/medstargenx/MEDSTARGENX/ml_api
Environment="PATH=/var/www/medstargenx/MEDSTARGENX/venv/bin"
ExecStart=/var/www/medstargenx/MEDSTARGENX/venv/bin/python /var/www/medstargenx/MEDSTARGENX/ml_api/cancer_api.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medstargenx-ml

[Install]
WantedBy=multi-user.target
